# .config/autostart/alacritty.desktop
[Desktop Entry]
Exec=alacritty
Icon=system-run
Name=alacritty
Path=
Terminal=True
Type=Application
