# .config/autostart/alacritty.desktop
[Desktop Entry]
Exec=alacritty
Icon=system-run
Name=Alacritty
Path=
Type=Application
