# .config/autostart/atril.desktop
[Desktop Entry]
Type=Application
Exec=atril
Hidden=false
Name[en_US]=Atril
Name=Atril
Comment[en_US]=
Comment=
X-MATE-Autostart-Delay=0
