# .config/autostart/st.desktop
[Desktop Entry]
Type=Application
Exec=st
Hidden=false
Name[en_US]=Suckless st
Name=Suckless st
Comment[en_US]=
Comment=
X-MATE-Autostart-Delay=0
